# userspace/root/etc/sv/dnsmasq/run
#!/bin/sh
exec 2>&1

CONFIG="/run/dnsmasq-usb0.conf"

mkdir -p /var/lib/misc

while [ ! -f "$CONFIG" ]; do
  sleep 1
done

exec dnsmasq --keep-in-foreground --conf-file="$CONFIG"
