# userspace/root/etc/sv/usb-hub-gpio/run
#!/bin/sh
# One-shot service - run once then sleep forever
/usr/comma/usb_hub_gpio.sh
exec sleep infinity
