# userspace/sv/dnsmasq/run
#!/bin/sh
exec 2>&1
mkdir -p /var/lib/misc
exec dnsmasq --keep-in-foreground --conf-dir=/etc/dnsmasq.d
