# userspace/root/etc/sv/avahi-alias-publish/run
#!/bin/sh
exec /usr/comma/avahi-publish-aliases.sh
